# programs/autostart/ledfx.desktop
[Desktop Entry]
Exec=sh -c "sleep 5 && ledfx"
Icon=
Name=LedFX
Path=
Terminal=false
Type=Application
